# scripts/label_studio_vps.service
[Unit]
Description=fable Label Studio review pack (80 images, loopback-only, auth required)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fable-label-studio
EnvironmentFile=/etc/fable-label-studio.env
# Resource bounds — VPS is shared with fable-dashboard (MemoryMax=900M)
MemoryMax=1400M
MemoryHigh=1100M
CPUQuota=150%
TasksMax=256
# Do not enable job executor paths; this unit is labeling-only
# Bind loopback only — public access is nginx TLS on 8081 → 127.0.0.1:8082
ExecStart=/opt/fable-label-studio/.venv/bin/label-studio start \
  --host 127.0.0.1 \
  --port 8082 \
  -b \
  --data-dir /opt/fable-label-studio/data
Restart=on-failure
RestartSec=5
# Harden process environment visibility somewhat (password still in EnvironmentFile only)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
